# root/etc/services.d/video-transcoder/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc python /usr/bin/transcoder.py